# rhinopi-maintenance/net-check/net-check.timer
[Unit]
Description=Timer for keeping internet connection active

[Timer]
OnCalendar=hourly
Unit=net-check.service

[Install]
WantedBy=timers.target
