# rhinopi-maintenance/git-sync/git-sync.service
[Unit]
Description=Service for keeping git repos in sync

[Service]
Type=oneshot
ExecStart=/bin/sh /root/scripts/git-sync/sync.sh