# mini-surv/systemd/upload/minisurv-mega-upload.timer
[Unit]
Description=Timer for minisurv mega upload

[Timer]
OnCalendar=*:0/3
Unit=minisurv-mega-upload.service

[Install]
WantedBy=timers.target