# mini-surv/systemd/update/minisurv-update.service
[Unit]
Description=Service for minisurv update

[Service]
Type=oneshot
ExecStart=/bin/sh /root/mini-surv/scripts/update.sh