# mini-surv/systemd/update/minisurv-update.timer
[Unit]
Description=Timer for minisurv update

[Timer]
OnCalendar=00/6:*
Unit=minisurv-update.service

[Install]
WantedBy=timers.target