# rhinopi-maintenance/tunnels/tunnel-restart.service
[Unit]
Description=Restarts tunnels

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart serveo-tunnel.service
